# ci/1020-docker-build/files/root/root/profile.d/10002.nvm.sh
. ~/.nvm/nvm.sh
export METEOR_ALLOW_SUPERUSER=true
